# config/gratia-xdmod.cron
41 * * * * root /usr/bin/gratia-xdmod -s 900

